# templates/default/uwsgi_init_script.erb
#!/bin/bash
#
# This script starts and stops "<%= @app_name %>"
#
# chkconfig:   2345 80 20
# description: "<%= @app_name %>" service
# processname: "<%= @app_name %>"

. /etc/rc.d/init.d/functions

DESC="<%=@app_name %>"

PROG=/usr/local/bin/uwsgi
DAEMON_OPTS="--ini <%= @ini_file_config %>"


PIDFILE=/var/run/<%= @app_name %>.pid
lockfile=/var/lock/subsys/<%= @app_name %>


start() {
    echo $"Starting $DESC: "
    # use daemon to start the service 
    daemon "$PROG $DAEMON_OPTS" --pidfile $PIDFILE
    retval=$?
    echo
     return $retval
}

stop() {
    echo $"Stopping $DESC: "
    killproc -p $PIDFILE -TERM
    retval=$?
    echo
    [ $retval -eq 0 ] && rm -f $lockfile
    return $retval
}

restart() {
    stop
    start
}

reload() {
    restart
}

force_reload() {
    restart
}

rh_status() {
    # run checks to determine if the service is running or use generic status
    status -p $PIDFILE
}

rh_status_q() {
    rh_status >/dev/null 2>&1
}


case "$1" in
    start)
        #Only if not running,  start
        rh_status_q && echo "<%= @app_name %> already started..." && exit 0
        $1
        ;;
    stop)
        # Only if service is running stopt
        rh_status_q || echo "<%= @app_name %> not started..." || exit 0
        $1
        ;;
    restart)
        $1
        ;;
    reload)
        rh_status_q || exit 7
        $1
        ;;
    force-reload)
        force_reload
        ;;
    status)
        rh_status
        ;;
    condrestart|try-restart)
        rh_status_q || exit 0
        restart
        ;;
    *)
        echo $"Usage: $0 {start|stop|status|restart|condrestart|try-restart|reload|force-reload}"
        exit 2
esac
    exit $?